# scripts/com.ecommerce.ingest-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ecommerce.ingest-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__PROJECT__/run_pipeline.py</string>
        <string>--watch</string>
        <string>--interval</string>
        <string>5</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT__</string>

    <key>StandardOutPath</key>
    <string>__PROJECT__/logs/watcher.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT__/logs/watcher-error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
